# 10.0/pkg/mariadb100.rc
#!/bin/ksh

daemon="${TRUEPREFIX}/mariadb-10.0/bin/mysqld_safe"

. /etc/rc.d/rc.subr

pexp="${TRUEPREFIX}/mariadb-10.0/libexec/mysqld .*"
rc_bg=YES
rc_reload=NO

rc_cmd $1
